# systemd/cloudflare-ddns.service
[Unit]
Description=Cloudflare Dynamic DNS Updater
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
# Update 'dynamic-dns-cloudflare' below if you used a different directory name
# %h expands to user's home directory
ExecStart=%h/.local/bin/uv run %h/dynamic-dns-cloudflare/cloudflare_ddns.py
WorkingDirectory=%h/dynamic-dns-cloudflare

# Logging goes to journald automatically
StandardOutput=journal
StandardError=journal

# Restart behavior handled by timer, not service
